# glassfish/templates/default/glassfish-init.d-script.erb
#!/bin/bash
#
# glassfish:          Startup script for Glassfish Application Server.
#
# chkconfig: 3 80 05
# description:      Startup script for domain1 of Glassfish Application Server.

GLASSFISH_HOME=<%= File.join(@node[:glassfish][:INSTALL_HOME],"glassfish") %>
export GLASSFISH_HOME

GLASSFISH_OWNER=<%= @node[:glassfish][:systemuser] %> ;
export GLASSFISH_OWNER

start() {
        echo -n "Starting Glassfish: "
        echo "Starting Glassfish at `date`" >> $GLASSFISH_HOME/domains/domain1/logs/startup.log
        su $GLASSFISH_OWNER -c "$GLASSFISH_HOME/bin/asadmin start-domain domain1" >> $GLASSFISH_HOME/domains/domain1/logs/startup.log
        sleep 2
        echo "done"
}

stop() {
        echo -n "Stopping Glassfish: "
        echo "Stopping Glassfish at `date`" >> $GLASSFISH_HOME/domains/domain1/logs/startup.log
        su $GLASSFISH_OWNER -c "$GLASSFISH_HOME/bin/asadmin stop-domain domain1" >> $GLASSFISH_HOME/domains/domain1/logs/startup.log
        echo "done"
}

# See how we were called.
case "$1" in
        start)
                start
                ;;
        stop)
                stop
                ;;
        restart)
                stop
                start
                ;;
        *)
                echo $"Usage: glassfish {start|stop|restart}"
                exit
esac
