# pen/templates/default/pen.upstart.conf.erb
author "Avishai Ish-Shalom <avishai@fewbytes.com>"
description "Upstart job for Pen load balancer cluster: <%= @name %>"

start on started networking

expect fork
respawn
respawn limit 5 15
exec /usr/bin/pen <%= @pen_options %> <%= @port %> <%= @pen_nodes.join(" ") %>
