# bluepill/templates/centos/bluepill_init.erb
#!/bin/sh
#
# Author: Jamie Winsor (<jwinsor@riotgames.com>)
#
# chkconfig: 345 99 1
# Provides: <%= @service_name %>
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6

BLUEPILL_BIN=<%= node['bluepill']['bin'] %>
BLUEPILL_CONFIG=<%= @config_file %>
SERVICE_NAME=<%= @service_name %>

case "$1" in
  start)
    echo "Loading bluepill configuration for $SERVICE_NAME "
    $BLUEPILL_BIN load $BLUEPILL_CONFIG
    ;;
  stop)
    $BLUEPILL_BIN $SERVICE_NAME stop
    $BLUEPILL_BIN $SERVICE_NAME quit
    ;;
  restart)
    $0 stop
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac
